# scripts/kizashi.cron
# Kizashi crontab (root 運用, /root/kizashi)。
# Install with: sudo crontab -e  (then paste the line below)
#
# cron has a minimal PATH, so use the ABSOLUTE path to uv.
# 抽出は cron ではなく常駐サービス kizashi-agent-worker (API課金ゼロ) で回す。
# 詳細は DEPLOY.md を参照。

# Collection every 3 hours (free, no API cost)
0 */3 * * * cd /root/kizashi && /root/.local/bin/uv run kizashi-daily --no-enrich >> /root/kizashi/run.log 2>&1
