# scripts/kizashi-agent-worker.service
[Unit]
# Kizashi 常駐エージェント抽出ワーカー (課金ゼロ / ログイン済み claude CLI)
# 収集 cron が貯めた未処理プールを、スコアの高い順に少しずつ抽出し続ける。
# 導入:
#   sudo cp scripts/kizashi-agent-worker.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now kizashi-agent-worker
#   journalctl -u kizashi-agent-worker -f      # ログ追尾
Description=Kizashi agent enrichment worker (subscription claude CLI, zero API cost)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/kizashi
# claude CLI のログイン情報 (~/.claude) を参照できるよう HOME を明示。
Environment=HOME=/root
Environment=KIZASHI_DB=/root/kizashi/kizashi.db
# バッチ5件→短休止。プールが空になると 15 分アイドルして収集 cron の新規を待つ。
ExecStart=/root/.local/bin/uv run kizashi-agent-worker --batch 5 --sleep 5 --idle-sleep 900 --db /root/kizashi/kizashi.db
Restart=always
RestartSec=30
# 常駐だがCPUは控えめに (抽出はI/O待ちが主)。
Nice=10

[Install]
WantedBy=multi-user.target
